# debian/dists/stable/termux/source/bashrc-d/data/data/com.termux/files/usr/etc/profile.d/qq
/home/ww/tur/debian/dists/stable/termux/source/bashrc-d/data/data/com.termux/files