# dists/stable/main/source/service-clashd/data/data/com.termux/files/usr/var/service/clashd/run
#!/data/data/com.termux/files/usr/bin/sh
exec clash -d $HOME/.config/clash/ > $PREFIX/tmp/.clashd 2>&1
